# templates/prune-restic@.timer.erb
[Unit]
Description=Schedule backup pruning for %i

[Timer]
OnBootSec=<%= @prunebootdelay %>
OnUnitActiveSec=<%= @prunefrequency %>

[Install]
WantedBy=multi-user.target

