# templates/restic@.timer.erb
[Unit]
Description=Schedule backups for %i

[Timer]
OnBootSec=<%= @bootdelay %>
OnUnitActiveSec=<%= @frequency %>

[Install]
WantedBy=multi-user.target

